# playbooks/vllm_template.service
[Unit]
Description=

[Service]
EnvironmentFile=
Type=simple
ExecStart=
User=

[Install]
WantedBy=multi-user.target
